# deploy/paperclip.service
[Unit]
Description=Paperclip Control Plane
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/paperclip
EnvironmentFile=/root/paperclip/deploy/.env
ExecStart=/usr/bin/pnpm --filter @paperclipai/server dev
Restart=always
RestartSec=5
TimeoutStopSec=30
KillSignal=SIGTERM
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
